# profile.d/pyenv.sh
#! /usr/bin/env sh
hash pyenv > /dev/null 2> /dev/null && eval "$(pyenv init --path)"
